# .codeploy/portable.fs/usr/lib/systemd/system/ops-nodejs@.socket
[Unit]
PartOf=ops-nodejs.service

[Socket]
ListenStream=%i
